# web_srv/bina_rag.service
[Unit]
Description=Bina RAG service

[Service]
ExecStart=/home/gwolf/start_rag.sh

[Install]
WantedBy=multi-user.target